# software/etc/webserver.supervisord.conf
# -*- conf -*-

[program:webserver]
command = python /opt/webserver/main.pyc
autostart=true
autorestart=true
numprocs=1
numprocs_start=8000
