# software/etc/webserver.supervisord.conf
# -*- conf -*-

[program:webserver]
command = python /opt/infomatic/software/webserver/main.py -p 80
autostart=true
autorestart=true
numprocs=1
numprocs_start=8000
